# deploy/linux/kanotls.service
[Unit]
Description=kanotls TLS+Noise tunnel service
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kanotls
Restart=always
RestartSec=5s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
